# cd/etc/inittab
::sysinit:/etc/init.d/rcS
tty1::askfirst:-/bin/hush
tty2::askfirst:-/bin/hush
tty3::askfirst:-/bin/hush
tty4::askfirst:-/bin/hush
tty5::askfirst:-/bin/hush
::restart:/bin/init
::ctrlaltdel:/bin/reboot
